# diskwatch.service.example
[Unit]
Description=DiskWatch disk usage monitoring webapp
After=network.target

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=__INSTALL_DIR__
Environment=DB_PATH=__INSTALL_DIR__/diskwatch.db
Environment=CONFIG_PATH=__INSTALL_DIR__/config.yaml
ExecStart=__INSTALL_DIR__/venv/bin/uvicorn server:app --host 127.0.0.1 --port 8070
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
